# .vscode-server/data/User/History/18067112/WOPT.service
[Unit]
Description=Honeypot Database Monitor
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=root
ExecStart=/root/mainnet-indexer/scripts/restore-honeypot.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
